# rootfs/usr/lib/opkg/info/adblock-fast.list
/etc/init.d/adblock-fast
/etc/config/adblock-fast
/etc/uci-defaults/90-adblock-fast
