# rootfs/usr/lib/opkg/info/minidlna.list
/usr/sbin/minidlnad
/etc/init.d/minidlna
/etc/config/minidlna
/etc/sysctl.d/30-minidlna.conf
